# etc/init.d/boot_inet_restriction
#!/bin/sh /etc/rc.common

START=18
STOP=85

start () {
	logger "Setting Internet restrictions based on boot switch position"
	/etc/scripts/boot_inet_restriction
}

stop () {}
